# deploy/systemd/lb-agent-web.service
[Unit]
Description=LB-AGENT web
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/lb-agent/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/lb-agent/.venv/bin/python -m uvicorn app.web.api:create_app --factory --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
